# deploy/skyclaw.service
# batabeto — Personal DevOps AI Agent
#
# INSTALL:
#   sudo cp deploy/skyclaw.service /etc/systemd/system/skyclaw.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable skyclaw
#   sudo systemctl start skyclaw
#
# LOGS:
#   sudo journalctl -fu skyclaw

[Unit]
Description=batabeto — Personal DevOps AI Agent
Documentation=https://github.com/nagisanzenin/skyclaw
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/.skyclaw

# Kill any stale process before starting — prevents 409 Telegram conflicts
ExecStartPre=/bin/bash -c 'pkill -f "/usr/local/bin/skyclaw" 2>/dev/null || true; sleep 1'

ExecStart=/usr/local/bin/skyclaw start

ExecStop=/bin/kill -TERM $MAINPID
TimeoutStopSec=30
KillMode=mixed

Restart=always
RestartSec=10
StartLimitIntervalSec=120
StartLimitBurst=5

EnvironmentFile=-/root/.skyclaw/.env
Environment=RUST_LOG=info
Environment=RUST_BACKTRACE=0
Environment=KUBECONFIG=/etc/rancher/k3s/k3s.yaml
LimitNOFILE=65535

StandardOutput=journal
StandardError=journal
SyslogIdentifier=batabeto

[Install]
WantedBy=multi-user.target
